# src/third_party/chromiumos-overlay/sys-boot/kdump-init/files/upstart/kdump-init.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "ChromeOS kdump init"
author          "chromium-os-dev@chromium.org"

start on started boot-services

oom score never

exec /usr/share/kdump/kdump-init.sh
