# src/platform2/virtual_file_provider/init/virtual-file-provider-cgroup.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chrome OS Virtual File Provider"
author          "chromium-os-dev@chromium.org"

start on starting system-services

script
  # Create freezer cgroup for virtual-file-provider. This allows us to freezer
  # processes in an order that prevents FUSE from making system suspend timeout.
  cgroup=/sys/fs/cgroup/freezer/virtual-file-provider
  mkdir -p "${cgroup}"
  chown power:power "${cgroup}/freezer.state"
  chown virtual-file-provider:virtual-file-provider "${cgroup}/cgroup.procs"
end script
