# src/platform2/flex_bluetooth/init/flex-bluetooth.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Run Reven-specific overrides for Floss Bluetooth"
author          "chromium-os-dev@chromium.org"

oom score -100

start on starting bluetoothd
task
exec /usr/bin/process_flex_bluetooth_overrides
