# src/platform2/bootid-logger/log-bootid-on-boot.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Record the boot time and boot id."
author          "chromium-os-dev@chromium.org"

oom score -100
# This limit is high to read boot_id.log.
limit as 100000000 unlimited

start on started syslog
task

exec /usr/sbin/bootid-logger
