# src/overlays/baseboard-trogdor/chromeos-base/chromeos-bsp-baseboard-trogdor/files/verify_fsg.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Load and verify modem FS image"
author        "chromium-os-dev@chromium.org"

# We want to start as soon as boot-services kick off so that we get to work
# ASAP.  The rmtfs task will wait until we're stopping before it starts.
start on starting boot-services
task

# This is just a script. It shouldn't be caught by the OOM killer.
# If killed, LTE might not work on this boot.
oom score -100

exec /usr/share/cros/init/verify_fsg.sh
