# src/platform2/trunks/placeholders/cr50-result.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS trunks pre-init event"
author          "chromium-os-dev@chromium.org"

start on started boot-services

oom score -100

exec logger -t "${UPSTART_JOB}" "cr50-result for non-GSC case performs no-ops."
