# src/platform2/timberslide/init/ecloggers.conf
# Copyright 2018 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "EC log concatenators"
author        "chromium-os-dev@chromium.org"

start on started system-services

task

script
  for path in /sys/kernel/debug/cros_*/console_log
  do
    start timberslide LOG_PATH="${path}"

    cros_variant="$(basename "$(dirname "${path}")")"
    if [ "${cros_variant}" = "cros_scp" ]; then
      start timberslide-watcher LOG_PATH="${path}"
    fi
  done
end script
