# src/third_party/chromiumos-overlay/media-libs/vpd-hdcpkey-install/files/udev/rules.d/99-rk3399-hdcp.rules
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
#
# Pass HDCP Key when we see the DP driver load.

ACTION=="add|change", DRIVER=="cdn-dp", RUN+="/lib/udev/hdcp_pass_key.sh %p"
