# src/overlays/overlay-sarien/chromeos-base/modemfwd-helpers/files/94-l850gl-gpio.rules
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

ACTION=="add", DEVPATH=="/devices/pci0000:00/INT34BB:00/gpiochip0/gpio/gpio311", \
    RUN+="/bin/chown modem:modem /sys/class/gpio/gpio311/direction /sys/class/gpio/gpio311/value"
