# src/third_party/chromiumos-overlay/chromeos-base/chromeos-activate-date/files/activate_date.service
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

[Unit]
Description=Update the Activation Date
PartOf=system-services.target
After=tlsdated.target

[Service]
ExecStart=/usr/share/cros/init/activate_date.sh
