# src/third_party/chromiumos-overlay/chromeos-base/chromeos-activate-date/files/activate_date.conf
# Copyright 2013 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Update the Activation Date"
author      "chromium-os-dev@chromium.org"

start on started tlsdated
stop on stopping system-services

exec /usr/share/cros/init/activate_date.sh
