# src/platform2/typecd/init/typecd.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the Chromium OS USB Type C daemon"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
expect fork
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

# Typecd *should* be able to recover from crashes (rebuild state when we start),
# so better to get OOM-killed than cause a memory panic.
oom score -100

# Let the daemon crash if it grows too much. "as" is "address space" (vm
# size). We expect a typical VM size of about 20MB for the daemon (currently)
# so set a limit for 5x that.
limit as 100000000 unlimited

exec minijail0 --config /usr/share/minijail/typecd.conf                      \
  -- /usr/bin/typecd                                                         \
