# src/overlays/chipset-cml/chromeos-base/power-governor-chipset-cml/files/99-set-governor-on-ac.rules
# If battery is charging, change the scaling_governor to performance.
# If battery is discharging, change the scaling_governor to powersave.
#
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/bin/cpufreq_config"
