# src/overlays/project-labstation/chromeos-base/chromeos-bsp-labstation/files/init/make-var-tmp-exec.conf
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Make the /tmp and /var dir executable"
author        "chromium-os-dev@chromium.org"

start on started system-services

script
  mount -o remount,exec "/tmp"
  mount -o remount,exec "/var"
end script
