# src/overlays/project-labstation/chromeos-base/chromeos-bsp-labstation/files/init/servod-usb-lock-file.conf
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Create servod usb lock file to coordinate usb detection"
author        "chromium-os-dev@chromium.org"

start on started system-services

script
  SERVOD_USB_LOCK_FILE="/var/lib/servod/lock_file"

  if [ ! -f "${SERVOD_USB_LOCK_FILE}" ] ; then
    touch "${SERVOD_USB_LOCK_FILE}"
    chmod 444 "${SERVOD_USB_LOCK_FILE}"
  fi
end script
