# src/third_party/chromiumos-overlay/sys-firmware/fwupd-storage/files/wilco/91-fwupdtool-storage-wdc.rules
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev rules to start firmware update for NVME SSD: WDC
SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x15b7", \
ATTR{firmware_rev}!="20240012", \
ATTR{model}=="PC SN520 NVMe*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=3099e006-2dd5-5285-80ec-2845d613dc53"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x15b7", \
ATTR{firmware_rev}!="11121012", \
ATTR{model}=="PC SN730 NVMe*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=10215a6e-5eb4-5fe4-a3d0-1ff4b3295840"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x15b7", \
ATTR{firmware_rev}!="21113012", \
ATTR{model}=="PC SN530 NVMe*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=84ebf8f2-1346-5712-acc7-c2feeed50e44"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x15b7", \
ATTR{firmware_rev}!="73103012", \
ATTR{model}=="PC SN740 NVMe*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=f8d46dc9-4450-5bf1-8f3d-46999571c4e6"
