# src/platform2/oobe_config/etc/init/ui_stop_waiter.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Block ui from stopping"
author          "chromium-os-dev@chromium.org"

# See oobe_config_save.conf for a description of these jobs.
start on stopping ui
stop on stopped oobe_config_save

# Process is not system critical but can not recover if killed.
oom score -400

# Wait until stop condition is met.
task
script
  while :; do sleep 1; done
end script
