# src/overlays/overlay-amd64-generic-koosh/chromeos-base/kexec-init/files/kexec.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Run the kexec at boot"
author          "cros-fuzzing@google.com"

start on starting startup

task

oom score -100

exec kexec -a -p /boot/vmlinuz --reuse-cmdline
