# src/platform2/wifi-testbed/init/wifi-testbed-init.conf
# Copyright 2013 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Set iptables policies and add regulatory rules for testbed AP"
author          "chromium-os-dev@chromium.org"

start on stopped iptables and stopped ip6tables

script
  # Disable the firewall.  This involves both adding rules and changing policies
  # so that DHCP traffic goes through.  Why this would require both is something
  # of a mystery.
  iptables -P INPUT ACCEPT
  iptables -P FORWARD ACCEPT
  iptables -P OUTPUT ACCEPT

  iptables -I INPUT -j ACCEPT
  iptables -I FORWARD -j ACCEPT
  iptables -I OUTPUT -j ACCEPT

  ip6tables -P INPUT ACCEPT
  ip6tables -P FORWARD ACCEPT
  ip6tables -P OUTPUT ACCEPT

  ip6tables -I INPUT -j ACCEPT
  ip6tables -I FORWARD -j ACCEPT
  ip6tables -I OUTPUT -j ACCEPT
end script
