# src/platform2/trunks/csme/tpm_tunneld.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS TPM tunnel daemon for CSME"
author          "chromium-os-dev@chromium.org"

start on started boot-services and started trunksd
stop on stopping boot-services

respawn

expect fork

exec tpm_tunneld
