# src/platform2/ippusb_bridge/init/printscan-debug-cleanup.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Turn off advanced printing/scanning debugging"
author          "chromium-os-dev@chromium.org"

# Starts on stopping ui so it runs on logout.
start on stopping ui
task

oom score never

# Turn off advanced logging for printing and scanning.
exec dbus-send --system --type=method_call --dest=org.chromium.debugd \
     /org/chromium/debugd org.chromium.debugd.PrintscanDebugSetCategories \
     uint32:0
