# src/third_party/chromiumos-overlay/chromeos-base/swap-init/files/init/swap.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Setup compressed memory swap"
author          "chromium-os-dev@chromium.org"

env MIN_LOW_MEMORY_MARGIN=

# For migration from swap.sh to swap_management, we make this config start
# after swap_start.conf, so that swap_start.conf can send dbus message to
# swap_management to start swap first. swap.sh will start swap alternately if
# swap_management fails.
start on stopped swap_start
task

exec /usr/share/cros/init/swap.sh start
