# src/third_party/chromiumos-overlay/chromeos-base/swap-init/files/init/swap.service
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

[Unit]
Description=Setup compressed memory swap
After=system-services.target

[Service]
ExecStart=/usr/share/cros/init/swap.sh start
