# src/third_party/chromiumos-overlay/chromeos-base/swap-init/files/init/teardown-swap.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Teardown swap on shutdown"
author          "chromium-os-dev@chromium.org"

# For migration from swap.sh to swap_management, we make this config start
# after swap_stop.conf, so that swap_stop.conf can send dbus message to
# swap_management to stop swap first. swap.sh will stop swap alternately if
# swap_management fails.
start on stopped swap_stop
task

exec /usr/share/cros/init/swap.sh stop
