# src/third_party/chromiumos-overlay/chromeos-base/selinux-policy/files/50-selinux.rules
SUBSYSTEM=="cpu", DRIVERS=="processor", DEVPATH=="/devices/system/cpu/*", RUN+="/sbin/restorecon -R /sys/devices/system/cpu"
