# src/overlays/baseboard-corsola/chromeos-base/lte_power_control/files/lte_power_on.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Init gpio power sequence for Fibocon LTE module"
author          "chromium-os-dev@chromium.org"

respawn
respawn limit 3 10

oom score -100

start on starting boot-services
task

exec lte_power_control on
