# src/platform2/attestation/pca_agent/server/pca_agentd.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS device PCA agent service."
author          "chromium-os-dev@chromium.org"

start on started boot-services and started shill and started syslog
stop on stopping boot-services

oom score -100

respawn

expect fork

script
  PROFILING_ARGS=""
  if [ -d "/mnt/stateful_partition/unencrypted/profraws" ]; then
    PROFILING_ARGS="-b /mnt/stateful_partition/unencrypted/profraws,,1"
  fi

  exec minijail0 --config /usr/share/minijail/pca_agentd.conf \
      ${PROFILING_ARGS} \
      -- /usr/sbin/pca_agentd
end script
