# src/overlays/overlay-geralt/chromeos-base/dfd_collector/files/dfd_collector.conf
# Copyright 2022 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Collect dfd-dump if hardware watchdog resets"
author          "chromium-os-dev@chromium.org"

oom score -100

start on starting boot-services
task

exec dfd_collector
