# src/platform2/imageloader/init/imageloader-shutdown.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "cleanup of component updater mounts"
author "chromium-os-dev@chromium.org"

start on stopping pre-shutdown
task

oom score never

# Cleanup all imageloader mounts so other filesystems can unmount cleanly.
exec /usr/sbin/imageloader --unmount_all
