# src/platform2/imageloader/init/imageloader-init.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "imageloader initializer"
author "chromium-os-dev@chromium.org"

start on starting system-services

oom score never

exec /usr/sbin/imageloader --init_only
