# src/platform2/chromeos-config/init/cros_configfs.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Abstract job to control the /run/chromeos-config mount"
author          "The ChromiumOS Authors <chromium-os-dev@chromium.org>"

start on starting udev

pre-start exec cros_config_setup
post-stop exec umount /run/chromeos-config/v1
