# release-R116-15509.B/src/platform2/login_manager/init/upstart/send-uptime-metrics.conf
# Copyright 2013 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report UMA statistics each time login prompt is shown"
author        "chromium-os-dev@chromium.org"

# This is triggered every time the login prompt is shown.
start on login-prompt-visible
stop on stopping system-services

# OK to fail. We'll get it next time.
oom score -100
task

exec /usr/share/cros/init/send-uptime-metrics
