# src/platform2/featured/share/featured.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS feature daemon"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping pre-shutdown
respawn

# This service enables kernel features. It can be killed,
# at the cost of some features not being enabled.
oom score -100

exec /sbin/featured
