# src/third_party/chromiumos-overlay/chromeos-base/hiberman/files/hiberman-teardown.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Tear down the hibernate device mapper devices and logical volumes when the user logs out"
author        "chromium-os-dev@chromium.org"

start on stopping ui

# The hiberman teardown command is short-lived, security critical and doesn't
# consume many resources.
oom score never

exec /usr/sbin/hiberman teardown-hiberimage
