# src/third_party/chromiumos-overlay/net-misc/rmtfs/files/77-rmtfs.rules
ACTION=="add", KERNEL=="qcom_rmtfs_mem?", GROUP=="rmtfs", MODE="0660", RUN+="/sbin/initctl emit --no-wait qcom-rmtfs-added"
ACTION=="add", KERNEL=="remoteproc?", DRIVERS=="qcom-q6v5-mss", RUN+="/sbin/initctl emit --no-wait qcom-modem-added"