# src/platform2/vm_tools/reference_vm/data/usr/local/lib/systemd/system/opt-google-cros\x2dcontainers.mount
[Unit]
DefaultDependencies=no
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What=LABEL=cros-vm-tools
Where=/opt/google/cros-containers
Options=ro
TimeoutSec=10

[Install]
WantedBy=local-fs.target
