# src/platform2/vm_tools/reference_vm/data/usr/local/lib/systemd/system/vshd.service
[Unit]
Description=vshd
Requires=opt-google-cros\x2dcontainers.mount
After=opt-google-cros\x2dcontainers.mount

[Service]
ExecStart=/opt/google/cros-containers/bin/vshd

[Install]
WantedBy=basic.target
