# src/overlays/overlay-zork/chromeos-base/chromeos-bsp-zork/files/90-xhci-quirks.rules
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

ACTION=="offline", SUBSYSTEM=="usb", ENV{ERROR}=="DEAD", RUN+="/lib/udev/xhci-restart.sh"
