# src/platform2/login_manager/init/upstart/login.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Mark that the user has logged in"
author          "chromium-os-dev@chromium.org"

start on start-user-session
task

script
  # The "logged-in" file is used by ACPI events to behave
  # appropriately.  The state directory is created earlier in
  # session startup.
  touch /run/state/logged-in
end script
