# src/platform2/login_manager/init/upstart/logout.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Mark that the user has logged out"
author          "chromium-os-dev@chromium.org"

# Use 'start on stopping' and 'task' to hold up restarting the
# session manager until we're done with marking the state change.
start on stopping ui
task

script
  # Indicate that we're no longer logged in
  rm -f /run/state/logged-in
end script
