# src/platform2/login_manager/init/upstart/send-uptime-metrics.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report UMA statistics each time login prompt is shown"
author        "chromium-os-dev@chromium.org"

# This is triggered every time the login prompt is shown.
start on login-prompt-visible

exec /usr/share/cros/init/send-uptime-metrics
