# src/platform2/login_manager/init/upstart/ui-collect-machine-info.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Collect machine information for session_manager"
author        "chromium-os-dev@chromium.org"

# This job collects machine information and places it in a file. It is
# then pushed to session_manager by ui-init-late.

# Start after the system is mostly up and running, which also implies
# that the VPD has been read and the udev database is fully initialized.
start on started system-services

exec /usr/share/cros/init/ui-collect-machine-info
