# src/platform2/login_manager/init/upstart/ui-respawn.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Reboot device upon repeated user interface crash."
author        "chromium-os-dev@chromium.org"

start on stopping ui

exec /usr/share/cros/init/ui-respawn
