# src/platform2/power_manager/init/upstart/cpufreq.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Configure cpufreq"
author          "chromium-os-dev@chromium.org"

# When performing normally, cpu_freq will exit ~100ms, but oom termination is
# not too bad.
oom score -100

start on started failsafe
task

exec cpufreq_config
