# src/platform2/power_manager/init/upstart/powerd.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the Chromium OS power daemon"
author        "chromium-os-dev@chromium.org"

# OK to kill if absolutely needed.
oom score -100

# Note, the 'start on' line is modified by chromeos-factoryinstall ebuild.
# Please do not change it without also changing that reference.
start on started boot-services
stop on stopping boot-services
expect fork
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

env LOG_DIR=/var/log/power_manager
env POWER_RUN_DIR=/run/power_manager/power
env PREFS_DIR=/var/lib/power_manager
env ROOT_RUN_DIR=/run/power_manager/root
env ROOT_SPOOL_DIR=/var/spool/power_manager/root
env VPD_CACHE_FILE=/mnt/stateful_partition/unencrypted/cache/vpd/filtered.txt
env MAX_NITS_PREF=internal_backlight_max_nits
env PM_TEST_DELAY_FILE=/sys/power/pm_test_delay

# Add modules here in --vmodule format for debugging.
env VMODULE_ARG=
env UPSTART_RUN=1

pre-start exec /usr/share/cros/init/powerd-pre-start.sh

exec minijail0 --config /usr/share/minijail/powerd.conf -- \
  /usr/bin/powerd                             \
    --log_dir=${LOG_DIR}                      \
    --run_dir=${POWER_RUN_DIR}                \
    --vmodule="${VMODULE_ARG}"                \
    >/var/log/powerd.out 2>&1
