# src/platform2/power_manager/init/upstart/report-power-metrics.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report resume metrics once per boot"
author        "chromium-os-dev@chromium.org"

# In normal operation, this script runs quickly, but is OK to kill if needed.
oom score -100

start on started system-services

# Report resume metrics after boot in case a shutdown occurred while
# the computer was suspended
exec send_metrics_on_resume -b
