# src/platform2/power_manager/init/upstart/thermal_zone.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Configure thermal zones"
author          "chromium-os-dev@chromium.org"

start on started system-services
task

# Unlikely to run out of memory, short lived process executed at boot.
oom score 100

# Arbitrary number given negligible memory usage.
limit as 10000000 unlimited

exec thermal_zone_config
