# src/third_party/chromiumos-overlay/chromeos-base/chromeos-cr50-scripts/files/cr50-disable-sleep.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
#
# This upstrat script just starts cr50-disable-sleep.sh, which monitors dbus
# for suspend events and then sends a command to the Cr50 to disable
# deep sleep during suspend.
#
# Keeping Cr50 awake during suspend consumes slightly moore power, but
# significantly improves resume time.
#
# This upstart script will only be installed when the
# `cr50_disable_sleep_in_suspend` USE flag is set

description     "Disable deep sleep on Cr50"
author          "chromium-os-dev@chromium.org"

start on started trunksd
stop on stopped trunksd

oom score -100
respawn

exec /usr/share/cros/cr50-disable-sleep.sh
