# src/third_party/chromiumos-overlay/chromeos-base/chromeos-base/files/udev-rules/99-uinput.rules
# Make the /dev/uinput device only accessible by the uinput group and notify
# services it has been made accessible.
ACTION=="add", SUBSYSTEM=="misc", KERNEL=="uinput", OWNER="root", GROUP="uinput", MODE="0660", RUN+="/sbin/initctl emit --no-wait uinput-added"
