# src/platform2/cronista/upstart/cronista.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start cronista a tamper evident log daemon."
author          "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

# Sacrifice before OOM panic.
oom score 0

# Run cronista with the default vsock address and port
exec minijail0 -u cronista -- /usr/bin/cronista -U vsock://vsock:3:5554
