# src/platform2/arc/container/myfiles/arc-myfiles-read.conf
# Copyright 2018 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Mount /home/chronos/user/MyFiles for ARC."
author        "chromium-os-dev@chromium.org"

start on started arc-myfiles
stop on stopping arc-myfiles

# Default score for Android system processes.
oom score -100

env SOURCE=/home/chronos/user/MyFiles
env DEST=/run/arc/media/MyFiles-read
# The read-only MyFiles should be achieved by allowing read (and prohibiting
# write) to group everybody. This is consistent with removable-read, and
# accessible by system apps (media process, DocumentsUI) and non-system apps
# with READ_EXTERNAL_STORAGE.
env UMASK=027
env ANDROID_ROOT_UID=0
env ANDROID_EVERYBODY_GID=9997

pre-start exec logger -t "${UPSTART_JOB}" "Pre-start ${UPSTART_JOB}"

# syslog-cat is used to redirect stdio from sdcard to the logging daemon.
# minijail is used in the mount-passthrough-jailed script.
exec /usr/sbin/syslog-cat --identifier="${UPSTART_JOB}" -- \
  mount-passthrough-jailed --source="${SOURCE}" --dest="${DEST}" \
    --fuse_umask="${UMASK}" --fuse_uid="${ANDROID_ROOT_UID}" \
    --fuse_gid="${ANDROID_EVERYBODY_GID}" --android_app_access_type="read"

post-stop exec logger -t "${UPSTART_JOB}" "Post-stop ${UPSTART_JOB}"
