# src/third_party/chromiumos-overlay/net-misc/qc-netmgr/files/qc-netmgr.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start qc-netmgr server"
author        "chromium-os-dev@chromium.org"

start on qrtr-service-added
stop on qrtr-service-removed

# This service keeps the modem data connection alive. It can be killed,
# at the cost of an interruption in LTE connectivity.
oom score -100
respawn

expect fork

# As per QC, we need to restart modemmanager everytime we start qc-netmgr.
pre-start script
  restart_modemmanager() {
    initctl restart modemmanager || true
  }

  restart_modemmanager
end script

# mmdata_mgr needs CAP_NET_ADMIN to open AF_QIPCRTR socket.
# TODO(evgreen): Put back proper sandboxing:
#   -S /usr/share/policy/qc-netmgr-seccomp.policy \

exec minijail0 --profile=minimalistic-mountns -inNlpvr --uts \
  -k "tmpfs,/sys,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC" \
  -b "/sys/devices/platform/soc@0/1e40000.ipa" \
  -c cap_net_admin=e -u nobody -g nobody \
  -- /usr/bin/mmdata_mgr
