# src/platform2/power_manager/udev/optional/93-powerd-wilco-ec-files.rules
DRIVER=="wilco-ec", \
    SUBSYSTEM=="platform", \
    RUN+="/bin/chown power:power \
        $sys$devpath/boot_on_ac \
        $sys$devpath/usb_charge", \
    RUN+="/bin/chmod ug+w \
        $sys$devpath/boot_on_ac \
        $sys$devpath/usb_charge"

DRIVER=="wilco-charge-schedule", \
    SUBSYSTEM=="platform", \
    RUN+="/bin/sh -c \"/bin/chown power:power \
        $sys$devpath/advanced_charging/* \
        $sys$devpath/peak_shift/*\"", \
    RUN+="/bin/sh -c \"/bin/chmod ug+w \
        $sys$devpath/advanced_charging/* \
        $sys$devpath/peak_shift/*\""

KERNEL=="wilco-charger", \
    SUBSYSTEM=="power_supply", \
    RUN+="/bin/chown power:power \
        $sys$devpath/charge_type \
        $sys$devpath/charge_control_start_threshold \
        $sys$devpath/charge_control_end_threshold", \
    RUN+="/bin/chmod ug+w \
        $sys$devpath/charge_type \
        $sys$devpath/charge_control_start_threshold \
        $sys$devpath/charge_control_end_threshold"
