# src/overlays/project-mobbase/chromeos-base/chromeos-firewall-init-mobbase/files/ip6tables.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Set iptables IPv6 policies and add rules"
author          "chromium-os-dev@chromium.org"

start on starting network-services
task

script
  # This script simply disables forwarding of traffic between interfaces
  # while leaving the firewall open for the services we plan to run as
  # part of mobbase projects.
  ip6tables -P FORWARD DROP
end script
