# src/platform2/init/upstart/ac-only.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Startup job to load AC only config settings"
author          "chromium-os-dev@chromium.org"

oom score never

start on started system-services

task

script
  psu="$(cros_config /hardware-properties psu-type)" || exit 0

  if [ "${psu}" = "AC_only" ]; then
    sysctl -p /etc/sysctl.d/*.conf.ac_only
  fi
end script
