# src/platform2/init/upstart/boot-alert-ready.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Abstract job that runs when boot-splash is done"
author          "chromium-os-dev@chromium.org"

# Not a real job; just propagates events.
oom score never

# The boot-alert-request signal is emitted by the chromeos-boot-alert
# script.  The initctl command there will block until the boot-alert-ready
# job has started, which won't happen until the boot-splash job has
# stopped.
start on boot-alert-request and stopped boot-splash
