# src/platform2/init/upstart/log-rotate.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Rotate the system logs"
author          "chromium-os-dev@chromium.org"

oom score -100

start on starting system-services
stop on stopping system-services

pre-start exec /usr/sbin/chromeos-cleanup-logs --no-periodic

exec periodic_scheduler --period=86400 --timeout=600 --task_name=cleanup_logs \
    -- /usr/sbin/chromeos-cleanup-logs
