# src/platform2/init/upstart/pstore.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "mount pstore for kernel buffer over reboots"
author          "chromium-os-dev@chromium.org"

# Unlikely to fail, but OK if it does (informational/crash diagnostics).
oom score -100

# encrypted-reboot-vault uses pstore-pmsg to store a transient key that is
# accessed during chromeos_startup. Mount pstore before beginning startup.
start on starting startup
task

env PSTORE_GROUP=pstore-access

script
  path="/sys/fs/pstore"
  mount -n -t pstore -o nodev,nosuid,noexec pstore "${path}"
  chmod 0750 "${path}"
  # TODO(ellyjones): restrict more tightly (crosbug.com/35122)
  chgrp "${PSTORE_GROUP}" "${path}"
end script
