# src/platform2/init/upstart/report-boot-complete.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Record that system boot is complete"
author        "chromium-os-dev@chromium.org"

# Non-critical metrics.
oom score -100

# This job is the official timestamp for boot time measurement.  It
# needs to run separately before 'system-services' because the amount
# of work that starts with 'system-services' could otherwise delay the
# 'bootstat' timestamp and make boot time measurements less reliable.

start on starting boot-complete
task

exec bootstat boot-complete
