# src/platform2/init/upstart/send-recovery-metrics.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report recovery metrics to UMA"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services

# OK to fail. We'll get it next time.
oom score -100
task

script
  # Recovery histograms file is written by chromeos-install during recovery.
  RECOVERY_HISTOGRAMS="/mnt/stateful_partition/.recovery_histograms"
  if [ -f "${RECOVERY_HISTOGRAMS}" ]; then
    # Replay events from the histogram file to send the metrics.
    metrics_client -R "${RECOVERY_HISTOGRAMS}"
    rm -f "${RECOVERY_HISTOGRAMS}" # croslint: disable:
  fi
end script
