# src/platform2/init/upstart/send-reclamation-metrics.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report full /var reclamation to UMA"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services

# OK to fail. We'll get it next time.
oom score -100
task

script
  # Reclamation file is written by chromeos_startup if /var got full and
  # had to be reclaimed.
  reclamation_file="/mnt/stateful_partition/.reclaim_full_var"

  # If the file is present, we need to notify UMA.
  if [ -f "${reclamation_file}" ]; then
    metrics_client -s "Startup.ReclaimFullVar" 1

    # Remove reclamation file to prevent duplicate reports.
    rm -f "${reclamation_file}" # croslint: disable:
  fi
end script
