# src/platform2/init/upstart/sysrq-init.conf
# Copyright 2019 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Determine whether to allow certain keyboard hotkey sequences"
author          "chromium-os-dev@chromium.org"

oom score -100

# We start in parallel with pre-startup as we don't need or want to block the
# critical code path, and we want to be robust if the rest of the system fails,
# and don't break the system if we fail.  This is a very uncommon flow and all
# other init scripts should be triggering before/after pre-startup instead.
start on startup

exec /sbin/cros_sysrq_init
