# src/platform2/init/upstart/udev.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the udev daemon"
author          "chromium-os-dev@chromium.org"

oom score never

# Start udev alongside chromeos_startup, since it only needs to
# write to /run and does very little else itself other than set
# up a few static nodes. The general "trigger"s for devices are in
# udev-trigger-* jobs.
start on stopped pre-startup
respawn

# Create static device nodes from modules.devname.
# This needs to complete before udev itself because udev's static_node rules
# expect the nodes to already exist so it can chown/chgrp them.
pre-start exec static_node_tool

expect fork
exec udevd --daemon
