# src/platform2/init/upstart/udev-boot.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Abstract job that runs when udev starts at boot."
author         "chromium-os-dev@chromium.org"

oom score never

# Some early services depend on udev, but don't need to be rerun if udev
# restarts/crashes.
start on started udev
# No 'stop' condition and no main process/script means this job runs forever.
