# src/platform2/init/upstart/uinput.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Load uinput module, required by bluetooth daemon"
author          "chromium-os-dev@chromium.org"

# Disable OOM killer as we must never fail.
oom score never

start on started system-services

exec modprobe uinput
