# src/third_party/chromiumos-overlay/chromeos-base/node_exporter/files/init/node_exporter.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run node_exporter on boot"
author        "chromium-os-dev@chromium.org"

start on started network-services
stop on stopping network-services

respawn
# If the job respawns 3 times in 10 seconds, stop trying.
respawn limit 3 10
# Do not respawn if the service is terminated on purpose.
normal exit 0
expect daemon

# OOM score for system metrics.
oom score -100

# minijail arguments:
# -u node_exporter: run as node_exporter user.
# -g node_exporter: run as node_exporter group.
# --uts: enter a new UTS namespace.
# -N: enter a new cgroup namespace.
# -n: set the no_new_privs bit.
# -T static: static mode to lockdown pre-exec.
# -i: exit immediately after fork(2). The jailed process will run in
#     the background.
exec minijail0 -u node_exporter -g node_exporter \
  --uts -N -n -T static -i \
  -- /usr/local/sbin/node_exporter \
  --web.listen-address=:9090 \
  --collector.textfile.directory=/var/lib/node_exporter \
  2>&1 | logger -t "${UPSTART_JOB}" &
