# src/third_party/chromiumos-overlay/sys-firmware/fwupd-storage/files/brya/91-fwupdtool-storage-samsung.rules
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev rules to start firmware update for NVME SSD: Samsung
SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x144d", \
ATTR{firmware_rev}!="26300039", \
ATTR{model}=="PM991a NVMe*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=f49677c5-4921-4733-81fa-69a09abbcc4d"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x144d", \
ATTR{firmware_rev}!="46303039", \
ATTR{model}=="PM9B1 NVMe Samsung 256GB*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=b55d114c-0ca9-4eeb-9f5b-706421e0c9c6"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x144d", \
ATTR{firmware_rev}!="46303039", \
ATTR{model}=="PM9B1 NVMe Samsung 512GB*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=b55d114c-0ca9-4eeb-9f5b-706421e0c9c6"
