# src/platform2/feedback/init/feedback_daemon.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Feedback collection daemon"
author          "chromium-os-dev@chromium.org"

# The feedback daemon is responsible for receiving and forwarding feedback
# reports.
start on starting system-services
stop on stopping system-services
respawn

exec /usr/bin/feedback_daemon
