# src/platform2/minios/ramfs/etc/init/minios.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Startup for MiniOS script"
author          "chromeos-core-services@google.com"

start on started dbus and started frecon and started shill
expect fork
respawn

oom score never

exec /usr/bin/minios
