# src/platform2/minios/ramfs/etc/init/syslog.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Syslog daemon"
author        "chromeos-core-services@google.com"

start on starting boot-services
expect fork
respawn

#TODO(b/181169861): miniOs: Run minijail with pivot root and bind mount
exec minijail0 --config /usr/share/minijail/syslog.conf \
     /usr/sbin/rsyslogd -n -f /etc/rsyslog.chromeos -i /tmp/rsyslogd.pid
