# src/platform2/minios/ramfs/etc/init/system-services.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Stub system-services script"
author          "chromeos-core-services@google.com"

start on started boot-services
