# src/platform2/minios/ramfs/etc/init/update-engine.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Job to start update engine"
author          "chromeos-core-services@google.com"

start on started dbus
stop on stopping dbus

expect fork
respawn
respawn limit 10 20

oom score never

exec update_engine
