# src/platform2/ml_core/opencl_caching/init/opencl-cacher.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "OpenCL shader compilation cacher"
author "chromium-os-dev@chromium.org"

start on start-user-session and stopped imageloader-init
stop on stopping ui

oom score -100

expect fork

# This service will load the ML processing shaders and run a
# tuning process to compile and optimize them for that specific
# device GPU. This happens on the GPU itself, and so the sandbox
# needs access to that device.

exec minijail0 --config /usr/share/minijail/opencl-cacher.conf \
    -- /usr/bin/opencl_cacher
