# src/platform2/init/upstart/test-init/encrypted_stateful/create-system-key.conf
# Copyright 2019 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Chrome OS post-system-key-setting cleanup"
author        "chromium-os-dev@chromium.org"

# This job handles delayed dumping of the system key creation log file into
# syslog, until the encrypted stateful is mounted, but before the log of
# mount-encrypted is copied into /var/log and syslog.
start on starting mount-encrypted

script
  SYSTEM_KEY_LOG_FILE="/run/create_system_key.log"

  # Dump the create-system-key logs in syslog to preserve them over reboots.
  logger -t create-system-key -f "${SYSTEM_KEY_LOG_FILE}" || true
  rm -f "${SYSTEM_KEY_LOG_FILE}"
end script
