# src/platform2/crash-reporter/init/crash-sender.timer
[Unit]
Description=Run Crash Sender hourly

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h
