# src/third_party/chromiumos-overlay/chromeos-base/infineon-firmware-updater/files/tpm-firmware-update.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Upgrades TPM firmware during first boot"
author         "chromium-os-dev@chromium.org"

start on starting boot-services

# This job must complete before tcsd and the UI get launched. The
# latter depend on "started boot-services", declaring as task requires
# the job to complete before boot-services reach "started".
task

exec /usr/share/cros/init/tpm-firmware-update.sh
