# src/third_party/chromiumos-overlay/chromeos-base/infineon-firmware-updater/files/tpm-firmware-check.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Determines whether a TPM firmware update is available"
author         "chromium-os-dev@chromium.org"

start on started system-services and started tcsd

task

exec /usr/share/cros/init/tpm-firmware-check.sh
