# src/overlays/overlay-hana/chromeos-base/chromeos-bsp-hana/files/99-mtk_drm_hdmi_load.rules
ACTION=="add", DRIVERS=="mediatek-drm-hdmi", RUN+="/usr/sbin/hdcp_pass_key.sh %p"
