# src/overlays/overlay-hana/chromeos-base/chromeos-bsp-hana/files/99-mwifiex-wowlan.rules
# Enable wake on wifi for the mwifiex driver
KERNEL=="mlan*", SUBSYSTEM=="net", ACTION=="add", RUN+="/usr/sbin/ethtool -s %k wol u"
