# src/third_party/chromiumos-overlay/net-wireless/iwlwifi_rescan/files/60-iwlwifi.rules
SUBSYSTEM=="pci", ACTION=="change", DRIVER=="iwlwifi", ENV{EVENT}=="INACCESSIBLE", RUN+="log-wifi-gone-metric.sh"
# Intel 7265, StP2.
SUBSYSTEM=="pci", ACTION=="remove", ENV{PCI_ID}=="8086:095A", RUN+="/sbin/initctl emit --no-wait pci-rescan-to-revive-wifi DEVICEPATH=$devpath"
SUBSYSTEM=="pci", ACTION=="remove", ENV{PCI_ID}=="8086:095B", RUN+="/sbin/initctl emit --no-wait pci-rescan-to-revive-wifi DEVICEPATH=$devpath"
# Intel 7260,  WP2.
SUBSYSTEM=="pci", ACTION=="remove", ENV{PCI_ID}=="8086:08B1", RUN+="/sbin/initctl emit --no-wait pci-rescan-to-revive-wifi DEVICEPATH=$devpath"
SUBSYSTEM=="pci", ACTION=="remove", ENV{PCI_ID}=="8086:08B2", RUN+="/sbin/initctl emit --no-wait pci-rescan-to-revive-wifi DEVICEPATH=$devpath"
