# src/platform2/metrics/init/metrics-daemon.service
[Unit]
Description=Metrics Daemon

[Service]
ExecStart=/usr/bin/metrics_daemon --nodaemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
