# src/platform2/metrics/init/structured_metrics.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Structured Metrics upstart file"
author          "chromium-os-dev@chromium.org"

oom score never

start on starting boot-services

# Increment reset counter. The counter integer is needed for the CrOS Events
# project to capture event sequences.
exec minijail0 -u metrics -g metrics /usr/bin/update_reset_counter
