# src/platform2/cryptohome/init/cryptohomed-client.conf
# Copyright 2012 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Stop existing cryptohomed clients"
author          "chromium-os-dev@chromium.org"

# Any long-running invocations of the client must first be terminated.
start on stopping cryptohomed
task

oom score -200

exec pkill -x -TERM cryptohome
