# src/platform2/cryptohome/init/init-homedirs.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Initializes the home directory"
author        "chromium-os-dev@chromium.org"

start on starting attestationd or starting cryptohomed

task

oom score 0

exec /usr/sbin/homedirs_initializer
