# src/overlays/overlay-volteer/chromeos-base/chromeos-bsp-volteer/files/udev/79-sof-dmic.rules
SUBSYSTEM=="pci", ATTR{vendor}=="0x8086", ATTR{device}=="0xa0c8", ATTR{class}=="0x040100", ATTRS{[dmi/id]board_name}=="Eldrid", RUN+="/sbin/modprobe snd_sof_pci tplg_path=intel/sof-tplg/pdm1-drceq"
