# src/platform2/cecservice/share/cecservice.conf
# Copyright 2018 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "cecservice"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.
oom score -100

expect fork

exec minijail0 --config /usr/share/minijail/cecservice.conf \
     -- /usr/sbin/cecservice
