# src/platform2/shadercached/init/shadercached.conf
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Start the ChromiumOS shader cache daemon"
author "chromeos-gaming-core@google.com"

start on started system-services and started vm_concierge
stop on stopping system-services and stopped vm_concierge

# Minijail forks off the desired process.
expect fork

# shadercached is a restartable daemon.
oom score -100

# If the job respawns 10 times in 10 seconds, stop trying.
respawn
respawn limit 10 10

exec minijail0 \
  --config /usr/share/minijail/shadercached.conf \
  /usr/bin/shadercached
