# src/platform2/spaced/init/spaced.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Disk usage daemon for ChromeOS"
author          "chromium-os-dev@chromium.org"

# Start once d-bus is ready.
start on started boot-services
stop on stopping boot-services

# Minijail forks off the desired process.
expect fork

# If the job respawns 10 times in 10 seconds, stop trying.
respawn
respawn limit 10 10

# Let the daemon crash if it grows too much.
limit as 30000000 unlimited
# Allow spaced to be killed and restarted. Consumers of data from spaced
# (momentarily) be unable to get correct data.
oom score -100

exec minijail0 --config /usr/share/minijail/spaced.conf \
  -- /usr/sbin/spaced
