# src/third_party/chromiumos-overlay/net-libs/libqrtr/files/qrtr-ns.conf
# Copyright 2018 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the qrtr-ns daemon"
author        "ejcaruso@chromium.org"

start on started network-services
stop on starting pre-shutdown

expect fork

pre-start exec qrtr-cfg 1

# qrtr-ns needs CAP_NET_ADMIN
exec minijail0 -inNlpvr --uts -S /usr/share/policy/qrtr-ns-seccomp.policy \
    -u qrtr -g qrtr -c 0x1000 \
    /usr/bin/qrtr-ns -f -s
