# src/platform2/diagnostics/cros_minidiag/upstart/send-boot-elog-metrics.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report once-per-boot elog-related UMA statistics"
author        "chromium-os-dev@chromium.org"

# Because the elogtool leverages SPI driver, wait for ui.conf ready to
# ensure there is no race condition between ui and this script which
# might result in increasing the boot time.
start on started system-services and started ui

# Killable for memory leaks.
oom score -100

exec cros-minidiag-tool --metrics_launch_count --metrics_test_report \
     --last_report --update_last_report
