# src/overlays/baseboard-herobrine/chromeos-base/modemfwd-helpers-herobrine/files/modemfwd-check.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Check if modemfwd has mounted FW"
author        "pholla@chromium.org"

start on started startup
task
oom score -100   # Short task at startup; don't oom kill

script
    # Check if modemfwd has mounted FW. If it hasn't, bypass modemfwd
    # to bring up cellular. Give modemfwd ~90s to download
    # or install DLCs.
    sleep 90;
    if [ ! -e /lib/firmware/qcom/sc*/modem/mba.mbn ]; then
            logger -p err -t "${UPSTART_JOB}" \
            "Modem FW not loaded by modemfwd. Bypassing modemfwd."

            /sbin/initctl start \
                          modemfwd-helpers \
                          main=/opt/google/modemfwd-firmware/modem-latest
    fi
end script
