# src/third_party/chromiumos-overlay/sys-firmware/fwupd-storage/files/hatch/91-fwupdtool-storage-liteon.rules
# Copyright 2022 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev rules to start firmware update for NVME SSD: Liteon
SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x14a4|0x1e95", \
ATTR{firmware_rev}!="22301116", \
ATTR{model}=="CL1-3D128-Q11*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=1cd7457a-27a9-4af9-912b-73af953d94b4"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x14a4|0x1e95", \
ATTR{firmware_rev}!="22301116", \
ATTR{model}=="CL1-3D256-Q11*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=1cd7457a-27a9-4af9-912b-73af953d94b4"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x14a4|0x1e95", \
ATTR{firmware_rev}!="22321116", \
ATTR{model}=="CL1-3D512-Q11*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=d8b197d0-0b03-4138-8ea6-8366a0cfde4a"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x14a4|0x1e95", \
ATTR{firmware_rev}!="25301111", \
ATTR{model}=="CL4-3D256-Q11*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=690e40e2-2141-4710-93a6-8a2ba6890c99"

SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x14a4|0x1e95", \
ATTR{firmware_rev}!="25301111", \
ATTR{model}=="CL4-3D512-Q11*", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=690e40e2-2141-4710-93a6-8a2ba6890c99"

# b/264579155: "CL1-3D256" contains its vendor name inside the model name.
SUBSYSTEM=="nvme", \
ATTRS{vendor}=="0x14a4", \
ATTR{firmware_rev}!="CR22002", \
ATTR{model}=="LITEON CL1-3D256", \
ACTION=="add", \
RUN+="/sbin/initctl emit fwupdtool-update PLUGIN=nvme GUID=c7ab5f4e-213d-52fe-b26e-7fdc0c2f3d1c"
