# src/platform2/init/upstart/test-init/early-failure.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Simple service that just crashes. Used to test handling of early crashes.

description "Start a process that crashes early"
author      "chromium-os-dev@chromium.org"

# This is extremely unlikely to run out of memory.
oom score -100

start on stopped crash-reporter-early-init and stopping pre-startup
task

exec timeout --signal=SIGSEGV 0.01s sleep 10000
