# src/platform2/init/upstart/test-init/factory.conf
# Copyright 2013 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Chrome OS factory startup stub"
author        "chromium-os-dev@chromium.org"

oom score never

# The start/stop timing should be the same (or earlier than) to ui.conf
start on started boot-services and stopped factory-init
stop on starting pre-shutdown

import GOOFY_ARGS
exec /usr/local/factory/init/startup main
