# src/overlays/project-mobbase/chromeos-base/chromeos-bsp-mobbase/files/init/mobbase-python-fix-init.conf
# Copyright 2014 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Fix broken python symlink for mobbase projects."
author        "chromium-os-dev@chromium.org"

start on starting boot-services

script
  PYTHON_LOCAL_LINK=/usr/local/bin/python2.7
  if [ ! -e "${PYTHON_LOCAL_LINK}" ]; then
    mkdir -p /usr/local/bin
    ln -s /usr/bin/python2.7 "${PYTHON_LOCAL_LINK}"
  fi
end script
