# src/overlays/overlay-dedede/chromeos-base/chromeos-bsp-dedede/files/common/dedede-force-wp.conf
# Copyright 2020 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Enable FW WP on Magolor/Maglia rev1 to avoid firmware incompatibilities."
author        "chromiumos-dev@chromium.org"

start on started system-services
task
oom score -100   # Short task at startup; don't oom kill

script
  set -e

  logit() {
    logger -t "${UPSTART_JOB}" "$@"
  }

  # Trigger only on magolor or maglia.
  model="$(cros_config / name)"
  if [ "${model}" != "maglia" -a "${model}" != "magolor" ] ; then
    logit "v2: Leaving WP alone for ${model}"
    exit 0
  fi

  # Trigger only for rev1.
  rev="$(ectool cbi get 0 | grep 'As binary: ' | sed 's/As binary: //')"
  if [ "${rev}" -ne 1 ]; then
    logit "Leaving WP alone for ${model} rev${rev}"
    exit 0
  fi

  # Only blow WP if the RO FW is up to date: 13606.52.0 or newer.
  minfw_maj=13606
  minfw_min=52
  fwid="$(crossystem ro_fwid)"
  fwnum_maj="$(echo ${fwid} | sed -n 's/Google_Magolor\.\([0-9]*\)\..*/\1/p')"
  fwnum_min="$(echo ${fwid} | sed -n 's/Google_Magolor\.[0-9]*\.\([0-9]*\)\..*/\1/p')"
  if [ "${fwnum_maj}" -lt "${minfw_maj}" ]; then
    logit "Not enabling WP for old FW ${fwnum_maj}"
    exit 0
  fi

  if [ "${fwnum_maj}" -eq "${minfw_maj}" -a "${fwnum_min}" -lt "${minfw_min}" ]; then
    logit "Not enabling WP for old FW ${fwnum_maj}.${fwnum_min}"
    exit 0
  fi

  if crossystem cros_debug?1 ; then
    logit "Not write-protecting debug image."
    exit 0
  fi

  if ! crossystem mainfw_type?normal ; then
    logit "Not write-protecting devmode image"
    exit 0
  fi

  logit "Enabling WP for ${model} rev${rev} ${fwnum_maj}.${fwnum_min}"
  futility flash --wp-enable
end script
