# src/third_party/chromiumos-overlay/chromeos-base/chromeos-docker/files/init/dockerd.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Run the dockerd on chromeos"
author         "chromium-os-dev@chromium.org"

start on started ui
stop on stopped ui

oom score 0

# putting docker data on unecrypted stateful partition allows for maiximum
# space for the docker data
env DOCKER_DATA_ROOT="/mnt/stateful_partition/docker"

exec dockerd --data-root "${DOCKER_DATA_ROOT}"
