# src/platform2/arc/vm/mojo_proxy/init/arcvm-server-proxy.conf
# Copyright 2018 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the server proxy for arcvm"
author        "chromium-os-dev@chromium.org"

start on starting arcvm-post-login-services
stop on stopping arcvm-post-login-services

oom score -100

exec minijail0 --config /usr/share/minijail/arcvm-server-proxy.conf \
  -- /usr/bin/arcvm_server_proxy /mnt

post-stop exec logger -t "${UPSTART_JOB}" "Post-stop arcvm-server-proxy"
