# src/overlays/baseboard-oak/chromeos-base/chromeos-bsp-baseboard-oak/files/platform-cpusets-mt8173-4_4.conf
# Copyright 2017 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "set up cpusets for Chrome on mt8173-oak as per gru"
author          "chromium-os-dev@chromium.org"

start on stopped cgroups
task

script
  # Non-urgent Chrome threads should go to little cores.
  echo 0-1 > /sys/fs/cgroup/cpuset/chrome/non-urgent/cpus

  # We use a schedtune cgroup to boost the deadline sensitive tasks.
  schedtune_dir="/sys/fs/cgroup/schedtune/chrome"
  # default boost parameters needed for consistent 60fps compositing
  echo 20 > "${schedtune_dir}/urgent/schedtune.boost"
  echo 1 > "${schedtune_dir}/urgent/schedtune.prefer_idle"

  # enable touch boosting
  echo 75 > /sys/module/cpu_boost/parameters/input_boost_freq_percent
end script
