# src/overlays/baseboard-oak/chromeos-base/chromeos-bsp-baseboard-oak/files/platform-cpusets-mt8173-3_18.conf
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "set up cpusets for Chrome"
author          "chromium-os-dev@chromium.org"

start on stopped cgroups
task

script
  # Non-urgent Chrome threads should go to little cores.
  echo 0-1 > /sys/fs/cgroup/cpuset/chrome/non-urgent/cpus
  # Urgent Chrome threads should go to big cores.
  echo 2-3 > /sys/fs/cgroup/cpuset/chrome/urgent/cpus
end script
