# src/platform2/init/upstart/halt/halt.conf
# Copyright 2011 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Halt (power down) the system"
author          "chromium-os-dev@chromium.org"

oom score never

start on runlevel 0
stop on runlevel [!0]

# SHUTDOWN_REASON is passed to initctl by powerd_setuid_helper.
# Export it here so that the pre-shutdown job can log it.
import SHUTDOWN_REASON
export SHUTDOWN_REASON

exec chromeos_shutdown
post-stop exec halt --force --poweroff
