# src/platform2/modemfwd/modemfwd.conf
# Copyright 2018 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Modem firmware updater daemon"
author          "ejcaruso@chromium.org"

start on starting system-services and stopped imageloader-init
stop on stopping system-services

# Same oom score level as e.g. shill
oom score -100

# Set to true to turn on debug mode on startup.
env DEBUG_MODE="false"
import DEBUG_MODE

script
  # Override logging level if sticky flag exists.
  if [ -f /var/cache/modem-utilities/log_modemfwd_debug ]; then
    DEBUG_MODE="true"
  fi
  # Default helper directory is /opt/google/modemfwd-helpers and if the
  # firmware directory isn't specified we load from the cellular component.
  exec minijail0 --config /usr/share/minijail/modemfwd.conf \
      /usr/bin/modemfwd \
      --firmware_directory=/opt/google/modemfwd-firmware \
      --debug_mode="${DEBUG_MODE}"

end script
