# src/third_party/chromiumos-overlay/chromeos-base/chromeos-machine-id-regen/files/cros-machine-id-regen-periodic.service
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

[Unit]
Description=Regenerate machine-id periodically
PartOf=system-services.target
After=system-services.target

[Service]
# Configures a maximum time for the service to run.
RuntimeMaxSec=60s
ExecStart=/usr/sbin/cros-machine-id-regen -r periodic -t 21600
