# src/third_party/chromiumos-overlay/chromeos-base/chromeos-machine-id-regen/files/cros-machine-id-regen-network.service
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

[Unit]
Description=Regenerate machine-id when disconnecting from a network

[Service]
ExecStart=/usr/sbin/cros-machine-id-regen -r network

[WantedBy]
Install=shill-disconnected.target
