# src/third_party/chromiumos-overlay/chromeos-base/chromeos-machine-id-regen/files/cros-machine-id-regen-network.conf
# Copyright 2015 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Regenerate machine-id when disconnecting from a network"
author          "chromium-os-dev@chromium.org"

start on shill-disconnected

task

exec /usr/sbin/cros-machine-id-regen -r network
