# src/third_party/chromiumos-overlay/chromeos-base/chromeos-machine-id-regen/files/cros-machine-id-regen-periodic.conf
# Copyright 2015 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Regenerate machine-id periodically"
author          "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services

exec periodic_scheduler --period=3600 --timeout=60 \
     --task_name=cros-machine-id-regen -- \
     /usr/sbin/cros-machine-id-regen -r periodic -t $((6 * 3600))
