# src/third_party/chromiumos-overlay/chromeos-base/chromeos-machine-id-regen/files/cros-machine-id-regen-periodic.timer
# Copyright 2016 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

[Unit]
Description=Regenerate machine-id timer
PartOf=system-services.target
After=system-services.target

[Timer]
# The service is run 1h after the timer is activated
# and from that moment on hourly.
OnActiveSec=1h
# Time between running each consecutive time.
OnUnitActiveSec=1h
Unit=cros-machine-id-regen-periodic.service

[WantedBy]
Install=system-services.target
